# ops/peppr-api-watchdog.service.example
[Unit]
Description=PepPro API watchdog (consecutive failures)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/peppr-api-watchdog.sh
